# picuntu-linuxroot-0.9/etc/init/alsa-restore.conf
# alsa-restore - restore sound card(s') mixer state(s)
#
# This task restores the sound card(s') mixer state(s) to
# /var/lib/alsa/asound.state on startup

description	"restore sound card(s') mixer state(s)"

start on runlevel [2345]

task

exec /usr/sbin/alsactl restore
