# picuntu-linuxroot-0.9/etc/init/alsa-store.conf
# alsa-store - save sound card(s') mixer state(s)
#
# This task saves the sound card(s') mixer state(s) to
# /var/lib/alsa/asound.state on reboot or shutdown

description	"save sound card(s') mixer state(s)"

start on runlevel [!2345]

task

exec /usr/sbin/alsactl store
