# picuntu-linuxroot-0.9/etc/init/mounted-dev.conf
# mounted-dev - Populate /dev filesystem
#
# Populates the /dev filesystem from /lib/udev/devices once the temporary
# filesystem mount is in place.

description	"Populate /dev filesystem"

start on mounted MOUNTPOINT=/dev
env MOUNTPOINT=/dev

task

script
    # With the /run transition, shm lives there now, but eglibc still looks in
    # /dev/shm. So create a symlink there
    [ -e /dev/shm ] || ln -s /run/shm /dev/shm
    
    if [ "${TYPE}" != "devtmpfs" ]
    then
	cd "${MOUNTPOINT}"
	/sbin/MAKEDEV std console fd ppp tun
    fi
end script
