# picuntu-linuxroot-0.9/var/lib/dpkg/info/cgroup-lite.list
/.
/bin
/bin/cgroups-mount
/bin/cgroups-umount
/usr
/usr/bin
/usr/share
/usr/share/doc
/usr/share/doc/cgroup-lite
/usr/share/doc/cgroup-lite/copyright
/usr/share/doc/cgroup-lite/changelog.gz
/etc
/etc/init.d
/etc/init
/etc/init/cgroup-lite.conf
/usr/bin/cgroups-mount
/usr/bin/cgroups-umount
/etc/init.d/cgroup-lite
