# picuntu-linuxroot-0.9/var/lib/dpkg/info/module-init-tools.list
/.
/lib
/lib/modules
/sbin
/sbin/depmod
/sbin/rmmod
/sbin/insmod
/sbin/modinfo
/sbin/modprobe
/bin
/bin/lsmod
/etc
/etc/modprobe.d
/etc/modprobe.d/blacklist-framebuffer.conf
/etc/modprobe.d/blacklist-firewire.conf
/etc/modprobe.d/blacklist.conf
/etc/modprobe.d/iwlwifi.conf
/etc/modprobe.d/blacklist-ath_pci.conf
/etc/modprobe.d/blacklist-rare-network.conf
/etc/modprobe.d/blacklist-watchdog.conf
/etc/init
/etc/init/module-init-tools.conf
/etc/depmod.d
/etc/depmod.d/ubuntu.conf
/etc/init.d
/usr
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/lsmod.8.gz
/usr/share/man/man8/modinfo.8.gz
/usr/share/man/man8/depmod.8.gz
/usr/share/man/man8/rmmod.8.gz
/usr/share/man/man8/insmod.8.gz
/usr/share/man/man8/modprobe.8.gz
/usr/share/man/man5
/usr/share/man/man5/modules.dep.5.gz
/usr/share/man/man5/modules.dep.bin.5.gz
/usr/share/man/man5/depmod.conf.5.gz
/usr/share/man/man5/modprobe.conf.5.gz
/usr/share/man/man5/modules.5.gz
/usr/share/man/man5/depmod.d.5.gz
/usr/share/man/fr
/usr/share/man/fr/man5
/usr/share/man/fr/man5/modules.5.gz
/usr/share/doc
/usr/share/doc/module-init-tools
/usr/share/doc/module-init-tools/changelog.Debian.gz
/usr/share/doc/module-init-tools/copyright
/usr/share/doc/module-init-tools/NEWS.gz
/sbin/lsmod
/etc/init.d/module-init-tools
