# picuntu-linuxroot-0.9/var/lib/dpkg/info/rfkill.list
/.
/etc
/etc/init
/etc/init/rfkill-restore.conf
/etc/init/rfkill-store.conf
/etc/init.d
/usr
/usr/sbin
/usr/sbin/rfkill
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/rfkill.8.gz
/usr/share/doc
/usr/share/doc/rfkill
/usr/share/doc/rfkill/changelog.Debian.gz
/usr/share/doc/rfkill/copyright
/usr/share/doc/rfkill/README
/etc/init.d/rfkill-restore
/etc/init.d/rfkill-store
